# private_dot_config/systemd/exact_user/xsession@.target
[Unit]
Description=X session on %I
Wants=xsettingsd@%i.service
Wants=xss-lock@%i.service
