# private_dot_config/systemd/exact_user/xss-lock@.service
[Unit]
Description=Manage X screen saver on %I

[Service]
Environment=DISPLAY=%I
ExecStart=%h/.local/bin/xss-lock start
Restart=on-failure
